# GPService/systemd/gpservice.service.in
[Unit]
Description=GlobalProtect openconnect DBus service

[Service]
Environment="LANG=en_US.utf8"
Type=dbus
BusName=com.yuezk.qt.GPService
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/gpservice

[Install]
WantedBy=multi-user.target
